# gruvvim/config/kmonad/kmonad.service
[Unit]
Description=Kmonad service

[Service]
ExecStart=kmonad /etc/systemd/kmonad.conf

[Install]
WantedBy=multi-user.target
